# scripts/com.fingerprint.ollama-tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.fingerprint.ollama-tunnel</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/autossh</string>
        <string>-M</string>
        <string>0</string>
        <string>-o</string>
        <string>ServerAliveInterval=30</string>
        <string>-o</string>
        <string>ServerAliveCountMax=3</string>
        <string>-N</string>
        <string>-R</string>
        <string>172.17.0.1:11434:localhost:11434</string>
        <string>vps</string>
    </array>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/com.fingerprint.ollama-tunnel.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/com.fingerprint.ollama-tunnel.stderr.log</string>
</dict>
</plist>
